# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root
pidfile=/var/run/supervisord.pid
logfile=/var/log/supervisor/supervisord.log
childlogdir=/var/log/supervisor/
loglevel=info
minfds=10000

[program:laravel-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600 --queue=high,default
directory=/var/www/html
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-queue-worker.log
stopwaitsecs=3600
stopasgroup=true
priority=999
autostart=true
autorestart=true
startsecs=1

[program:laravel-reverb-server]
process_name=%(program_name)s
command=php /var/www/html/artisan reverb:start
directory=/var/www/html
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-reverb-server.log
stopwaitsecs=30
stopasgroup=true
priority=998
autostart=true
autorestart=true
startsecs=3

[program:laravel-scheduler]
process_name=%(program_name)s
command=php /var/www/html/artisan schedule:work
directory=/var/www/html
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-scheduler.log
stopwaitsecs=30
stopasgroup=true
priority=997
autostart=true
autorestart=true
startsecs=1

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[inet_http_server]
port=127.0.0.1:9001
